# debian/pt-display.postinst
#!/bin/bash -e

echo "Unmasking pt-scale-desktop-ui service..."
systemctl unmask pt-display-scale-desktop-ui

echo "Enabling pt-scale-desktop-ui service..."
systemctl enable pt-display-scale-desktop-ui

echo "Starting pt-scale-desktop-ui service..."
systemctl start pt-display-scale-desktop-ui

exit 0
