# debian/pt-os-mods.prerm
#!/bin/bash

systemctl enable rpi-eeprom-update.service
