# debian/pt-os-mods.prerm
#!/bin/bash

systemctl enable rpi-eeprom-update.service

sed -i "s|EEPROM_CONFIG_HOOK=\"/usr/lib/rpi-eeprom-update/pt-hook.sh\"|d" /etc/default/rpi-eeprom-update

# TODO: Uncomment out EEPROM_CONFIG_HOOK in /etc/default/rpi-eeprom-update if it exists
